# cardcopy@.service
[Unit]
Description=Copy Photos

[Service]
Type=simple
ExecStart=/usr/local/bin/cardcopy %I

